# etc/systemd/system/fwupd-refresh.service.d/override.conf
[Service]
ExecStart=/usr/bin/fwupdmgr update